# supvan-printer-app.service
[Unit]
Description=Supvan T50 Pro Printer Application
After=bluetooth.target dbus.socket

[Service]
Type=simple
ExecStart=/usr/bin/supvan-printer-app server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
